# systemd/kilo-mcp-ncp-validator.service
[Unit]
Description=Kilo MCP: NCP Validator (@ :8217)
After=network-online.target

[Service]
ExecStart=%h/.npm-global/bin/mcp-proxy --port 8217 -- node %h/.config/kilo/mcp-servers/ncp-validator/dist/index.js
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
