# systemd/kilo-mcp-ralph.service
[Unit]
Description=Kilo MCP: Ralph PRD execution (@ :8210)
After=network-online.target

[Service]
ExecStart=%h/.npm-global/bin/mcp-proxy --port 8210 -- node %h/.config/kilo/mcp-servers/ralph/dist/index.js
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
