# systemd/kilo-mcp-hooks.service
[Unit]
Description=Kilo MCP: Hooks & Channels (@ :8212)
After=network-online.target

[Service]
ExecStart=%h/.npm-global/bin/mcp-proxy --port 8212 -- node %h/.config/kilo/mcp-servers/hooks/dist/index.js
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
