# systemd/kilo-mcp-score-engine.service
[Unit]
Description=Kilo MCP: Score Engine (@ :8219)
After=network-online.target

[Service]
ExecStart=%h/.npm-global/bin/mcp-proxy --port 8219 -- node %h/.config/kilo/mcp-servers/score-engine/dist/index.js
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
