# systemd/kilo-mcp-constitution.service
[Unit]
Description=Kilo MCP: Constitution code quality (@ :8216)
After=network-online.target

[Service]
ExecStart=%h/.npm-global/bin/mcp-proxy --port 8216 -- node %h/.config/kilo/mcp-servers/constitution/dist/index.js
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
